# recipes-graphics/wayland/weston-init/weston.service.d/display-ready.conf
[Unit]
After=dev-dri-card0.device
Wants=dev-dri-card0.device
